# mac_server/Services/sentinelapi/com.hyperbrain.sentinelapi.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hyperbrain.sentinelapi</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Update path after first build -->
        <string>/Users/danielc/mac_server/Services/sentinelapi/.build/release/SentinelAPI</string>
        <string>serve</string>
        <string>--env</string>
        <string>production</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <!-- AWS credentials for the event-sentinel-api-sqs IAM user are read from the macOS
             keychain at startup (service com.hyperbrain.sentinelapi) — NEVER put them here. -->
        <key>AWS_REGION</key>
        <string>us-east-1</string>
        <!-- Full URL of sync-events.fifo (outbound Apple -> Core), provisioned in #26 -->
        <key>SYNC_EVENTS_QUEUE_URL</key>
        <string>REPLACE_WITH_SYNC_EVENTS_QUEUE_URL</string>
        <!-- Full URL of apple-commands.fifo (inbound Core -> Apple), provisioned in HU-09c #37 -->
        <key>APPLE_COMMANDS_QUEUE_URL</key>
        <string>REPLACE_WITH_APPLE_COMMANDS_QUEUE_URL</string>
        <!-- Full URL of apple-commands-results.fifo (outbound results, ADR-010) -->
        <key>APPLE_COMMANDS_RESULTS_QUEUE_URL</key>
        <string>REPLACE_WITH_APPLE_COMMANDS_RESULTS_QUEUE_URL</string>
        <!-- Tailscale hostname of Mac Mini — never use IP (can change) -->
        <key>SENTINEL_API_BASE_URL</key>
        <string>REPLACE_WITH_TAILSCALE_HOSTNAME</string>
        <!-- Optional: force the REST bind address. By default the binary binds to the
             Tailscale interface (100.64.0.0/10) and degrades to loopback when it is down
             (RNF-09) — never pass "serve --hostname 0.0.0.0". -->
        <!-- <key>SENTINEL_HOSTNAME</key><string>100.x.y.z</string> -->
    </dict>

    <key>KeepAlive</key>
    <true/>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/danielc/Library/Logs/sentinelapi/sentinelapi.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/danielc/Library/Logs/sentinelapi/sentinelapi.stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/danielc/mac_server/Services/sentinelapi</string>
</dict>
</plist>
